# services/personal.service
[Unit]
Description=personal

[Service]
ExecStart=/usr/bin/ls

[Install]
WantedBy=multi-user.target
